# scripts/crontab/debian-upgrade.sh
#!/usr/bin/env sh

apt update
apt upgrade
apt autoremove
