# .root/etc/profile.d/amd.sh
#!/bin/sh

export LIBVA_DRIVER_NAME=radeonsi
export VDPAU_DRIVER=radeonsi
