# scripts/arch/pkgs/llm-meta/llm-meta.install
post_install() {
    systemctl enable --now ollama
}
