# scripts/arch/pkgs/llm-meta/llm-meta.install
post_install() {
    systemctl enable --now ollama
    printf "==> Edit the systemd config like this:\nsystemctl edit ollama.service\n[Service]\nEnvironment=\"OLLAMA_ORIGINS=https://ollama.sainnhe.dev\"\n"
}

pre_remove() {
    systemctl disable --now ollama
}
