# scripts/arch/pkgs/amd-meta/amd-meta.install
post_install() {
    systemctl enable --now lactd
    echo "==> Remember to install env file in /etc/profile.d"
}
