# .root/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://127.0.0.1:17080/"
Environment="HTTPS_PROXY=http://127.0.0.1:17080/"
