# scripts/arch/pkgs/cli-meta/cli-meta.install
post_install() {
    systemctl enable --now clash-meta
    systemctl enable --now podman
    systemctl enable --now sshd
}
