# scripts/arch/pkgs/hooks-meta/hooks-meta.install
post_install() {
    systemctl daemon-reload
    systemctl enable --now linux-modules-cleanup
}
