# scripts/arch/pkgs/kde-meta/kde-meta.install
post_install() {
    systemctl enable --now sddm
}
