# arch/packages/apps-common-meta/apps-common-meta.install
post_install() {
    systemctl enable --now bluetooth.service
}

pre_remove() {
    systemctl disable bluetooth.service
}
